# script/docker/supervisor-sidekiq.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=true
directory=/opt/alaveteli

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:redis-server]
command=redis-server

[program:sidekiq]
command=bundle exec sidekiq -L log/sidekiq.log -c 2

[program:clockwork]
command=bundle exec clockwork config/clock.rb

[program:sidekiq-web]
command=bundle exec thin -R sidekiq.ru -p 9292 start
